# src/boxes/pccc-box/root/etc/init.d/mount-fstab
#!/sbin/openrc-run

description="mount all fstab entries"

depend() {
	need net
}

start() {
	mount -a
	eend $?
}
