# src/boxes/linux-box/root/etc/init.d/ntpd
#!/sbin/openrc-run

description="update the time"

depend() {
	need net
}

start() {
	ntpd -d -q -n -p de.pool.ntp.org
	eend $?
}
