# src/boxes/docker-box/root/etc/init.d/update-ca-certificates
#!/sbin/openrc-run

description="update-ca-certificates"


start() {
	update-ca-certificates &> /dev/null
	eend 0
}
